# packages/cowpower/src/etc/udev/rules.d/80-cowpower.rules
SUBSYSTEM=="power_supply", ATTRS{online}=="1", ATTRS{type}=="Mains", ACTION=="change", RUN+="/usr/local/bin/cowpower"